# openpower/package/petitboot/66-add-sg-module.rules
# load modules to scsi disks, if they aren't in kernel
SUBSYSTEM=="scsi_device", ACTION=="add", RUN+="/sbin/modprobe sg"
