# openpower/overlay/etc/udev/rules.d/66-add-sg-module.rules
# load modules to scsi disks, if they aren't in kernel
SUBSYSTEM=="scsi", ENV{DEVTYPE}=="scsi_device", TEST!="[module/sg]", RUN+="/sbin/modprobe sg"
